# files/etc/uci-defaults/99_my_startup_enable
#!/bin/sh
/etc/init.d/my_startup_script enable
